# ansible/roles/i3-desktop/desktop/files/rc.local
#!/bin/sh

/bin/sh -c 'echo "XHC" > /proc/acpi/wakeup'